# playbooks/roles/worker-streaming/templates/nab3a-pipe.conf
[program:pipe 1]
command=php bin/nab3a pipe 1 -vv
user=debian
directory=/srv/kal3a
