# playbooks/roles/worker-streaming/templates/nab3a-pipe.conf
[program:pipe]
command=php bin/nab3a pipe 1 -vvv
user=debian
directory=/srv/kal3a
